# 0.10.3/yao-wms/crontab/crontab.conf
# 计划任务

## 数据统计
1 0 * * * /data/cloud/bin/stat.sh

## 数据备份
1 */1 * * * /data/cloud/bin/backup.sh

## 每天凌晨 5:30 重启服务器
30 5 * * * /sbin/shutdown -r +5

